# pi-camera-broadcaster.service
[Unit]
Description=Raspberry Pi Camera WebRTC Broadcaster for Stuffed Animal War
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=jaemzware
Group=jaemzware
WorkingDirectory=/home/jaemzware/stuffedanimalwar
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/jaemzware/stuffedanimalwar/venv/bin/python /home/jaemzware/stuffedanimalwar/pi_camera_broadcaster.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jaemzware/stuffedanimalwar

[Install]
WantedBy=multi-user.target
