# pisetup/camera-autostart.service
[Unit]
Description=Camera Autostart for StuffedAnimalWar
# Wait for graphical session and both application servers
After=graphical-session.target
After=stuffedanimalwar.service
After=analogarchive.service
Requires=stuffedanimalwar.service
Wants=analogarchive.service

[Service]
Type=simple
User=jaemzware
WorkingDirectory=/home/jaemzware/stuffedanimalwar
# Give the servers and X session time to fully initialize
ExecStartPre=/bin/sleep 10
# Format: node camera-autostart.js <domain> <room> [cameraName]
# This joins https://<domain>/<room>camera
# Change these values for your setup:
ExecStart=/usr/bin/node camera-autostart.js localhost:55556 frontdoor "Pi Camera"
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
# Puppeteer needs a display - use the default X display
Environment="DISPLAY=:0"
# Use system Chromium instead of downloading
Environment="USE_SYSTEM_CHROMIUM=true"

[Install]
WantedBy=graphical.target
